# test/inet-http-server-without-password.conf
[inet_http_server]
file = /tmp/supervisor.sock
port = 9001

[supervisord]
logfile = /tmp/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:cat]
command=/bin/cat
