# test/unix-http-server-with-password.conf
[unix_http_server]
file = /tmp/supervisor.sock
username = user
password = password

[supervisord]
logfile = /tmp/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:cat]
command=/bin/cat
